# config/access.conf
# Upstart Job to execute the Spring Boot application on startup
description "access"
start on startup
task
exec /vagrant/access/config/spring-start.sh
